# scheduler/crontab
# supercronic crontab (run by the `scheduler` service).
# Fields: min hour dom mon dow  command

# Scrape recent weather every day at 06:15 UTC.
15 6 * * *   portal weather fetch-recent --days 3

# Recompute the derived flow series nightly (cheap; derived data).
30 2 * * *   portal flow rebuild

# Nightly pg_dump (custom format) into /backups, keep ~14 days.
0 3 * * *    pg_dump -h db -U "$POSTGRES_USER" -d "$POSTGRES_DB" -Fc -f "/backups/skalna_$(date -u +\%Y\%m\%d).dump" && find /backups -name 'skalna_*.dump' -mtime +14 -delete
